# lab-data/security-imx93-frdm/userland/systemd/training-systemd-test.service
[Unit]
Description=Security training systemd demo

[Service]
ExecStart=/usr/local/bin/systemd-test
Type=oneshot
